# etc/update-motd.d/01-altherneum
#!/bin/bash

clear

[ -r /etc/lsb-release ] && . /etc/lsb-release

if [ -z "$DISTRIB_DESCRIPTION" ] && [ -x /usr/bin/lsb_release ]; then
    DISTRIB_DESCRIPTION=$(lsb_release -s -d)
fi

printf "рҹҗ§\tOS\t\t$(uname -o)\n"
printf "рҹ–ІпёҸ\tDistribution\t$DISTRIB_DESCRIPTION\n"
printf "рҹ“ҹ\tVersion\t\t$(uname -r)\t$(uname -m)\n"

printf "\n"

printf "вҢӣ\tUptime\t\t$(uptime -p)\n"
printf "рҹ“…\tDepuis\t\t$(uptime -s)\n"

printf "\n"

freeh=$(free -h)
freem=$(free -m)
meminfo=$(cat /proc/meminfo)
availmem="$(awk -v arcsize="0" '/^MemAvailable:/ { print $2 + 0 }' <<< "${meminfo}")"
usedmem="$(awk -v availmem="${availmem}" '/^MemTotal:/ { print $2 - availmem }' <<< "${meminfo}")"
totalmem="$(awk '/^MemTotal:/ { print $2 }' <<< "${meminfo}")"
used="$(numfmt --round=down --from-unit=1024 --to=iec <<< "${usedmem}")"
avail="$(numfmt --round=down --from-unit=1024 --to=iec <<< "${availmem}")"
total="$(numfmt --round=down --from-unit=1024 --to=iec <<< "${totalmem}")"
percentage=$(((usedmem*100 / totalmem*100) /100))

printf "рҹ’»\tCPU\t\t$(uptime | cut -d ',' -f 4-6 | sed -e 's/^[ \t]*//')\n"
printf "рҹ§ \tRAM\t\t$used / $total : $percentage%%\n"
printf "рҹ’ҫ\tSDA\n"
df -hl /dev/sda* | grep -v udev

printf "\n"

printf "рҹ’ҫ\tHГҙte\t\t$(hostname)\n"

printf "\n"

printf "рҹ“¶\tVPS\t\t$(hostname -I)\n"

printf "\n"

printf "рҹ‘Җ\tUtilisateurs\t$(who | sort -u | wc -l)\n"
printf "вһЎпёҸ\tListe :\n$(who)\n"

printf "\n"

printf "/var/log/wtmp (last) :"
printf "рҹ‘Җ\tDerniГЁres connexions\tрҹҶ—\n"
last | head

printf "\n"

printf "/var/log/btmp (lastb) :"
printf "рҹ‘Җ\tDerniГЁres connexions\tрҹ”ҙ\n"
lastb | head

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\tSSHD\n"
fail2ban-client status sshd

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\tCode-Server\n"
fail2ban-client status code-server

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\tKeeweb\n"
fail2ban-client status keeweb

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\tForgejo\n"
fail2ban-client status forgejo

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\t404\n"
fail2ban-client status 404

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\t.env\n"
fail2ban-client status env

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\t400\n"
fail2ban-client status 400

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\t401\n"
fail2ban-client status 401

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\t405\n"
fail2ban-client status 405

printf "\n"

printf "рҹҡЁ\tFail2Ban\tрҹ‘®\t403\n"
fail2ban-client status 403

printf "\n"

printf "рҹҡ§\tServices\t\tвҡЎ\t\tрҹ’Ө\n"
printf "вҳҺпёҸ\tDiscord\t\t\t$(systemctl is-enabled discord)\t\t$(systemctl check discord)\n"

printf "\n"

printf "рҹ“©\tsendmail\t\t$(systemctl is-enabled sendmail)\t\t$(systemctl check sendmail)\n"
printf "рҹ”җ\topendkim\t\t$(systemctl is-enabled opendkim)\t\t$(systemctl check opendkim)\n"

printf "\n"

printf "рҹҢҗ\tApache2\t\t\t$(systemctl is-enabled apache2)\t\t$(systemctl check apache2)\n"
printf "рҹ’ҝ\tMySQL\t\t\t$(systemctl is-enabled mysql)\t\t$(systemctl check mysql)\n"

printf "\n"

printf "рҹ“қ\tCode-Server\t\t$(systemctl is-enabled code-server@vscode)\t\t$(systemctl check code-server@vscode)\n"

printf "\n"

printf "рҹҗӢ\tdocker\t\t\t$(systemctl is-enabled docker)\t\t$(systemctl check docker)\n"

printf "\n"

printf "рҹ—ғпёҸ\tforgejo\t\t\t$(systemctl is-enabled forgejo)\t\t$(systemctl check forgejo)\n"
printf "рҹӣ пёҸ\tforgejo-runner\t\t$(systemctl is-enabled forgejo-runner)\t\t$(systemctl check forgejo-runner)\n"

printf "\n"

printf "рҹҡ«\tFail2Ban\t\t$(systemctl is-enabled fail2ban)\t\t$(systemctl check fail2ban)\n"

printf "\n"

printf "рҹ’Һ\tminecraft\t\t$(systemctl is-enabled minecraft)\t$(systemctl check minecraft)\n"

printf "\n"

printf "вҳ пёҸ\tminecraft@anarchie\t$(systemctl is-enabled minecraft@anarchie)\t$(systemctl check minecraft@anarchie)\n"
printf "рҹҸҳпёҸ\tminecraft@crГ©atif\t$(systemctl is-enabled minecraft@creatif)\t$(systemctl check minecraft@creatif)\n"
printf "рҹ‘®\tminecraft@opprison\t$(systemctl is-enabled minecraft@opprison)\t$(systemctl check minecraft@opprison)\n"
printf "вҡ”пёҸ\tminecraft@rpg\t\t$(systemctl is-enabled minecraft@rpg)\t$(systemctl check minecraft@rpg)\n"
printf "рҹҡ©\tminecraft@hub\t\t$(systemctl is-enabled minecraft@hub)\t\t$(systemctl check minecraft@hub)\n"
printf "рҹҡЎ\tminecraft@proxy\t\t$(systemctl is-enabled minecraft@proxy)\t\t$(systemctl check minecraft@proxy)\n"
printf "рҹҺҮ\tminecraft@skyblock\t$(systemctl is-enabled minecraft@skyblock)\t\t$(systemctl check minecraft@skyblock)\n"
printf "рҹҚ—\tminecraft@survie\t$(systemctl is-enabled minecraft@survie)\t\t$(systemctl check minecraft@survie)\n"

printf "\n"

printf "рҹҡ§\tminecraft@Test\t\t$(systemctl is-enabled minecraft@test)\t\t$(systemctl check minecraft@test)\n"

printf "\n"